# wpadremap.service.in
[Unit]
Description=Drawing tablet pad remapper
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=@libexecdir@/wpadremap
Restart=on-failure
